# docker/supervisord/queue.conf
[program:queue]
directory=/app
command=php artisan queue:work
autostart=true
autorestart=true
stdout_logfile=/app/storage/logs/queue.log
