# docker/supervisord/schefule.conf
[program:schedule]
directory=/app
command=php artisan schedule:work
autostart=true
autorestart=true
stdout_logfile=/app/storage/logs/schedule.log
